# seerr/rootfs/etc/services.d/seerr-ha-actions/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting Seerr Home Assistant actions helper..."
exec node /etc/seerr-ha-actions/server.js
